# recipes-rdkb/sysint-broadband/files/nvram.mount
[Unit]
Before=local-fs.target

[Mount]
Where=/nvram
What=/dev/disk/by-label/nvram1

[Install]
WantedBy=local-fs.target
