# recipes-ccsp/ccsp/ccsp-common-library/checkturriswifisupport.service
[Unit]
Description=Check WIFI support on Turris Omnia device

[Service]
Type=forking
ExecStart=/bin/sh -c '/usr/ccsp/wifi/checkwifi.sh'
