# recipes-ccsp/ccsp/ccsp-common-library/checkturriswifisupport.path
[Unit]
Description=Check Turris Wifi Support Initialized Path

[Path]
PathExists=/tmp/pam_initialized
Unit=checkturriswifisupport.service

[Install]
WantedBy=multi-user.target
