# recipes-ccsp/ccsp/ccsp-common-library/turriswifiinitialized.path
[Unit]
Description=TurrisOmnia WIFI Initialized Path

[Path]
PathExists=/tmp/wifi_driver_initialized
Unit=ccspwifiagent.service

[Install]
WantedBy=multi-user.target
