# recipes-ccsp/ccsp/ccsp-common-library/wifiinitialized.path
[Unit]
Description= WiFi Initialized Path

[Path]
PathExists=/tmp/wifi_initialized
Unit=wifi-initialized.target

[Install]
WantedBy=multi-user.target
