# recipes-ccsp/ccsp/ccsp-common-library/wifi-initialized.target
[Unit]
Description=WiFi Initialization Complete
Documentation=man:systemd.special(7)
