# recipes-connectivity/wpa-supplicant/files/wpa_supplicant-global.service
[Unit]
Description=WPA supplicant daemon
After=CcspPandMSsp.service
StartLimitIntervalSec=120

[Service]
Type=simple
PIDFile=/var/run/wpa_supplicant-global.pid
ExecStart=/usr/sbin/wpa_supplicant -g/var/run/wpa_supplicant-global -P /var/run/wpa_supplicant-global.pid 
Restart=always

[Install]
WantedBy=multi-user.target

